# rootfs/etc/s6-overlay/s6-rc.d/n8n/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: n8n
# Runs the n8n service
# ==============================================================================

declare webhook_url

chmod +x "$0"

bashio::log.info "Starting n8n service..."

webhook_url=$(bashio::config 'webhook_url')
bashio::log.info "webhook_url: ${webhook_url}"

# Charger les variables d’environnement depuis /data/n8n_env
source /data/n8n_env

exec n8n start