# rootfs/etc/s6-overlay/s6-rc.d/debug-shell/up
/etc/s6-overlay/s6-rc.d/debug-shell/run
