# rootfs/etc/s6-overlay/s6-rc.d/fetch-n8n-key/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: n8n
# Fetches and exports n8n encryption key if not specified
# ==============================================================================

declare encryption_key
declare generated_key

chmod +x "$0"

bashio::log.info "Fetching n8n encryption key..."

encryption_key=$(bashio::config 'encryption_key')

if bashio::config.has_value 'encryption_key'; then
  bashio::log.info "Clé de chiffrement valide trouvée dans les options: ${encryption_key}"
else
  bashio::log.info "Aucune clé de chiffrement valide spécifiée, attente de génération par n8n..."
  declare -i timeout=120
  declare -i count=0
  while [ ! -f /root/.n8n/config ] && [ "${count}" -lt "${timeout}" ]; do
    sleep 1
    count=$((count + 1))
    bashio::log.debug "Attente de /root/.n8n/config... (${count}/${timeout} secondes)"
  done
  if [ -f /root/.n8n/config ]; then
    generated_key=$(jq -r '.encryptionKey // empty' /root/.n8n/config)
    if bashio::var.has_value "${generated_key}"; then
      bashio::config 'encryption_key' "${generated_key}"
      bashio::log.info "Clé générée par n8n et exportée dans la config de l'addon: ${generated_key}"
    else
      bashio::log.warning "Fichier /root/.n8n/config trouvé, mais aucune clé présente."
    fi
  else
    bashio::log.warning "Timeout atteint (${timeout} secondes) ou échec: /root/.n8n/config non créé."
  fi
fi

exit 0